# worker/supervisord.ini
[program:mapworker]
directory = /home/dotcloud/current/
command = php workers.php
environment = QUEUE=*
stopsignal = QUIT
stderr_logfile = /var/log/supervisor/%(program_name)s_error-%(process_num)s.log
stdout_logfile = /var/log/supervisor/%(program_name)s-%(process_num)s.log
numprocs = 1
process_name = "%(program_name)s-%(process_num)s"